# applets/point_of_sale/tools/posbox/overwrite_after_init/etc/cron.daily/tele
#!/bin/sh

wget 'http://localhost:9000/hw_drivers/check_certificate' >/dev/null